# cron/rclone-box-bisync.cron
*/5 * * * * rclone bisync -v box:BoxSync ~/BoxSync --log-file ~/.box-rclone-bisync.log